# .config/upstart/mutate.conf
description "mutate"

start on desktop-start
stop on desktop-end

respawn

exec mutate
