# .config/upstart/redshift.conf
description "redshift redness maker"

start on desktop-start
stop on desktop-end

respawn

# lat:lon is for SF
exec redshift-gtk -l37.7833:-122.4167 -t5500:3400
