# .config/upstart/zulip.conf
description "zulip chat"

start on desktop-start
stop on desktop-end

respawn

exec zulip --site https://ice.integration.grandrounds.com
